# src/mod/endpoints/mod_gsmopen/gsmlib/gsmlib-1.10-patched-13ubuntu/contrib/gsm-utils.cron.d
# /etc/cron.d/gsm-utils: crontab fragment for gsm-utils

*/5 * * * *	root	if [ -x /usr/bin/gsmsmsrequeue ]; then /usr/bin/gsmsmsrequeue; fi
